# centos7-supervisor/demo.ini
[program:ftp]
command=/usr/bin/python3 -m http.server